# whisper-dictation.service
[Unit]
Description=Whisper.cpp transcription server

[Service]
ExecStart=/home/nat/tools/whisper.cpp/build/bin/whisper-server \
    -m /home/nat/tools/whisper.cpp/models/ggml-large-v3.bin \
    --host 127.0.0.1 \
    --port 8178 \
    -l en \
    --convert \
    --tmp-dir /tmp
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
