# whisper-dictation.service
[Unit]
Description=Whisper.cpp transcription server

[Service]
# Adjust paths below to match your whisper.cpp install location.
ExecStart=%h/tools/whisper.cpp/build/bin/whisper-server \
    -m %h/tools/whisper.cpp/models/ggml-large-v3.bin \
    --host localhost \
    --port 8178 \
    -l en \
    --convert \
    --tmp-dir /tmp
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
